# deploy_tools/gunicorn.template.service
[Unit]
Description=Gunicorn server for DOMAIN

[Service]
Restart=on-failure
User=kowajanski
WorkingDirectory=/home/kowajanski/sites/DOMAIN
EnvironmentFile=/home/kowajanski/sites/DOMAIN/.env

ExecStart=/home/kowajanski/sites/DOMAIN/virtualenv/bin/gunicorn \
    --bind unix:/tmp/DOMAIN.socket \
    superlists.wsgi:application

[Install]
WantedBy=multi-user.target